# deploy/lmbsc-web.service
[Unit]
Description=LMBSC Hospitality Platform — Next.js frontend
After=network.target lmbsc-api.service
Wants=lmbsc-api.service

[Service]
Type=exec
User=aurel
Group=aurel
WorkingDirectory=/home/aurel/mita/apps/web/.next/standalone
ExecStart=/usr/bin/node server.js
Environment=PORT=3001
Environment=HOSTNAME=127.0.0.1
Restart=always
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lmbsc-web

EnvironmentFile=-/home/aurel/mita/.env.production

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/aurel/mita

[Install]
WantedBy=multi-user.target
