# deploy/lmbsc-api.service
[Unit]
Description=LMBSC Hospitality Platform — FastAPI backend
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=exec
User=aurel
Group=aurel
WorkingDirectory=/home/aurel/mita/apps/api
ExecStart=/home/aurel/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 4
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lmbsc-api

# Environment file — copy from deploy/.env.production.example and fill in secrets
EnvironmentFile=/home/aurel/mita/.env.production

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/aurel/mita

[Install]
WantedBy=multi-user.target
